# hardwaremgr/hardwaremgr.service
[Unit]
Description=Hardware Manager Service
After=network.target

[Service]
ExecStart=/usr/local/bin/hardwaremgr
Restart=always
User=pi
WorkingDirectory=/path/to/your/project # Replace with the actual path to your project

[Install]
WantedBy=multi-user.target